# deploy/lets-care-deploy.service
# /etc/systemd/system/lets-care-deploy.service
# One-shot: check master and deploy if changed. Driven by the .timer.
[Unit]
Description=Poll master and deploy lets-care-portugal
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=letscarevpn
Group=letscarevpn
WorkingDirectory=/opt/lets-care-portugal
ExecStart=/opt/lets-care-portugal/deploy/deploy.sh
# systemd will not start a second run while one is active -> no overlap/locking needed.
TimeoutStartSec=900
